# restic.cron
# Cronfile for running Restic
# Place in /etc/cron.d/restic
# Fredrik Boulund 2021
15 4 * * * root /root/restic-backup/restic_backup.sh
